# applications/AdaptiveApplication1/proc1.service
[Unit]
Description=Adaptive Application 2 process 1
[Service]
Type=forking
PIDFile=/var/run/proc1.pid