# applications/MachineStateManager/msm.service
[Unit]
Description=Machine state manager
[Service]
Type=forking
TimeoutStartSec=infinity
