# applications/AdaptiveApplication2/proc3.service
[Unit]
Description=Adaptive Application 2 process 3
[Service]
Type=forking