# applications/AdaptiveApplication2/proc2.service
[Unit]
Description=Adaptive Application 2 process 2
[Service]
Type=forking
PIDFile=/var/run/proc2.pid